# crontab
0 */5 * * * /usr/bin/python3 /home/biodit/nodered-fixer/fixer.py
